# config/daemon/finai-event-monitor.service
[Unit]
Description=FinAI Research Workflow — Event-Driven Research Trigger Monitor
Documentation=https://github.com/csmar432/FinAI-Research-Workflow
After=network.target

[Service]
# Adjust paths to match your installation.
# `setup-daemon.sh linux` will replace __PROJECT_ROOT__ with your actual path
# before installing to /etc/systemd/system/.
WorkingDirectory=__PROJECT_ROOT__
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=__PROJECT_ROOT__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=__PROJECT_ROOT__/.env

# Run as a non-root user (recommended: create a 'finai' user)
# User=finai
# Group=finai

ExecStart=__PROJECT_ROOT__/.venv/bin/python3 \
    __PROJECT_ROOT__/scripts/event_monitor.py \
    --daemon \
    --auto-trigger \
    --macro-scheduler \
    --log-file __PROJECT_ROOT__/logs/monitor.log \
    --pid-file __PROJECT_ROOT__/data/event_monitor.pid \
    --pipeline-type research_report \
    --language zh

# Restart on crash (but not on normal exit)
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=append:__PROJECT_ROOT__/logs/monitor-stdout.log
StandardError=append:__PROJECT_ROOT__/logs/monitor-stderr.log

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__PROJECT_ROOT__/data
ReadWritePaths=__PROJECT_ROOT__/logs
ReadWritePaths=__PROJECT_ROOT__/output
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
